# files/captainhook.template.service
[Unit]
Description=captainhook

[Service]
ExecStart=bash ${DIR}/start-captainhook.sh

[Install]
WantedBy=default.target
